# dm.service
[Unit]
Description=Display Manager
After=systemd-user-sessions.service

[Service]
ExecStart=dm.py

[Install]
Alias=dm.service
